# shawtybot.service
[Unit]
Description=SHAWTY VPN Telegram Bot
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/shawtybot
ExecStart=/path/to/shawtybot/venv/bin/python /path/to/shawtybot/bot.py
Restart=always
RestartSec=10
Environment=PATH=/path/to/shawtybot/venv/bin
StandardOutput=append:/path/to/shawtybot/bot.log
StandardError=append:/path/to/shawtybot/bot.log

[Install]
WantedBy=multi-user.target
